# config/autostart/pa-applet.desktop
[Desktop Entry]
Exec=/usr/bin/pa-applet
Name=pa-applet
Comment=PulseAudio control applet
Type=Application
